# tools/telemetry_dashboard/systemd/chronon3d-telemetry-prune.timer
[Unit]
Description=Chronon3D Telemetry — Weekly prune timer (Sunday 03:00)

[Timer]
OnCalendar=Sun *-*-* 03:00:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
