# tools/telemetry_dashboard/systemd/chronon3d-dashboard-flask.service
[Unit]
Description=Chronon3D Dashboard — Flask Backend (API + WebSocket)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Pyt/Chronon3d
Environment=CHRONON3D_DASHBOARD_PASSWORD=chronon3d_admin
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/venv/bin/python3 %h/Pyt/Chronon3d/tools/start_dashboard_shim.py 8000
Restart=always
RestartSec=3
StandardOutput=append:%h/.chronon3d/logs/flask_backend.log
StandardError=append:%h/.chronon3d/logs/flask_backend.log

[Install]
WantedBy=default.target
